# overlay/etc/cron.d/espocrm
* * * * * www-data cd /var/www/espocrm; /usr/bin/php -f cron.php > /dev/null 2>&1
